# config/service/go.service.example
# Go application systemd service template
# Placeholders will be replaced by deploy script

[Unit]
Description={APP_NAME} Go Application
After=network.target

[Service]
User={USER}
Group={GROUP}
WorkingDirectory={APP_PATH}
ExecStart={APP_PATH}/app
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier={APP_NAME}
EnvironmentFile=-{APP_PATH}/.env
Environment=GOCACHE=/var/cache/go-build
Environment=GOMODCACHE=/var/cache/go-mod
Environment=GOPATH=/var/cache/go

# Optional: Limit resources
LimitNOFILE=65535
PrivateTmp=true
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
